# shell-script_prototype/_shell_test/loop_create.service
[Unit]
Description=Create Loop Device

[Service]
ExecStart=/usr/local/bin/loop.sh
Type=oneshot

[Install]
WantedBy=default.target
